# Scripts Raspberry/crontab
@reboot /home/pi/startup.sh
*/5 * * * * /home/pi/synch-cloud.py